# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=python /opt/qiandao/web.py

[program:worker]
command=python /opt/qiandao/worker.py
